# itamae/roles/wire/templates/var/lib/machines/overlay/etc/systemd/system/open-ddns.timer
[Unit]
Description=open-ddns

[Timer]
OnUnitActiveSec=5min
OnStartupSec=5min
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
