# itamae/roles/wire/templates/etc/systemd/system/systemd-nspawn@overlay.service.d/dropin.conf
[Service]
ExecStart=
ExecStart=/usr/bin/systemd-nspawn --quiet --keep-unit --boot --link-journal=try-guest --network-veth -U --settings=override --machine=%i --load-credential=network.wireguard.private.default:%d/network.wireguard.private.default
LoadCredentialEncrypted=network.wireguard.private.default:/etc/wire.key
